# common/etc/udev/rules.d/automount.rules
ACTION=="add", ENV{DEVTYPE}=="partition", RUN+="domount %N"
